# deet.timer
[Unit]
Description='Sets the interval to run deet pacman parser'

[Timer]
OnUnitActiveSec=1hr
OnBootSec=5min

[Install]
WantedBy=multi-user.target
